# gpu-server/grow-gpu.service
[Unit]
Description=Grow GPU Server (FastAPI)
After=network.target

[Service]
Type=simple
User=growapi
Group=growapi
WorkingDirectory=/home/growapi/app/aiseed-tool/gpu-server
ExecStart=/home/growapi/app/aiseed-tool/gpu-server/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

# 環境変数（.env は config.py が読み込むので最小限）
Environment=HOME=/home/growapi
Environment=PATH=/home/growapi/.local/bin:/usr/local/bin:/usr/bin:/bin

# ログ
StandardOutput=journal
StandardError=journal
SyslogIdentifier=grow-gpu

# セキュリティ
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/growapi/app/aiseed-tool/gpu-server

[Install]
WantedBy=multi-user.target
